# s6/proxy/run
#!/usr/bin/with-contenv sh
/usr/bin/tinyproxy -d
