# s6/proxy/run
#!/usr/bin/with-contenv sh
/usr/bin/mitmdump -p 3000
