# s6/app/finish
#!/usr/bin/env sh
echo >&2 "service2 exited. code=${1}"
exec s6-svscanctl -t /var/run/s6/services

