# s6/app/run
#!/usr/bin/with-contenv sh
HTTPS_PROXY=http://localhost:3000 /server
