# s6/app/run
#!/usr/bin/with-contenv sh
HTTP_PROXY=http://localhost:8888 /server
